# local_tel_api/telegram-bot-api.service
[Unit]
Description=Telegram Bot API Local Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu
Environment=PATH=/home/ubuntu/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=LD_LIBRARY_PATH=/home/ubuntu/.local/lib:/usr/local/lib:/usr/lib/x86_64-linux-gnu:/lib/x86_64-linux-gnu
ExecStart=/home/ubuntu/.local/bin/telegram-bot-api \
    --api-id=YOUR_API_ID \
    --api-hash=YOUR_API_HASH \
    --local \
    --http-port=8081 \
    --dir=/home/ubuntu/telegram-bot-api-data \
    --verbosity=2 \
    --max-webhook-connections=100000 \
    --max-connections=100000

# 重启策略
Restart=always
RestartSec=10

# 安全设置（放宽一些限制）
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=false
ProtectHome=false
ReadWritePaths=/home/ubuntu/telegram-bot-api-data

# 日志设置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-bot-api

[Install]
WantedBy=multi-user.target 